# mcp.service
[Unit]
Description=mcp

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/home/morgan/mcp

[Install]
WantedBy=multi-user.target
